# qic_install_tool/shell/sso
#chkconfig: 345 92 8
#author:liujianli@gtadata.com
#布署的时候改成对应的项目根目录就OK
. /etc/rc.d/init.d/functions
PROJ_NAME=sso
BASE_DIRECTORY=/var/myApp/QIC/sso
PID_FILE=$BASE_DIRECTORY/server.pid
#LOG_FILE=$DIRECTORY/logs/system.out
#play脚本位置
DEAMON=//usr/local/play1.2.x_cust/play
#java安装目录
export JAVA_HOME=/usr/local/jdk1.6.0_33/
export LANG="en_US.UTF-8"
RET_VAL=0;
function start()
{
  if [ -f "$PID_FILE" ]
  then
         failure
         echo "$PROJ_NAME服务运行中，请先关闭"
  else
        # echo "开始启动。。。"
         cd $BASE_DIRECTORY && $DEAMON start>/dev/null
         success
         echo "$PROJ_NAME服务启动成功"

  fi
}
function stop()
{
  if [ -f "$PID_FILE" ]
  then
        cd $BASE_DIRECTORY && $DEAMON stop>/dev/null;
        sleep 3
        if [ -f "$PID_FILE" ]
           then
           #play 意外中止不会删除pid文件
           rm -f $PID_FILE
        fi
        success
        echo "$PROJ_NAME停止成功"
  else
       warn 
       echo "$PROJ_NAME没启动"
  fi

}
function status()
{
   if [ -f "$PID_FILE" ]
    then
         pid=$(cat $PID_FILE)
         if [ ! -d "/proc/$pid" ]
         then
           echo "$1未运行"
           return 0
         fi
         $(getMode)
         mode=$?
         if [ $mode -eq 1 ];
         then
             success
             echo "$1运行中. PID为:$pid,运行模式:生产模式"
         else
             success
             echo "$1运行中. PID为:$pid,运行模式:开发模式"
         fi
         #cat $BASE_DIRECTORY$1/server.pid
        # echo ""
    else
         failure
         echo "$PROJ_NAME未运行"
    fi


}
# 0:开发模式 1 生产模式
function getMode()
{


       mode=$(sed -n "/^application.mode=dev/p" $BASE_DIRECTORY/conf/application.conf)

       if [ ${#mode} -gt 0 ]
       then
           return 0
       else
           #echo "test"
           return 1
       fi
       #mode_in_chinese=$(( "$mode" = "application.mode=dev"?"开发模式":"生产模式" ))

      #echo $mode_in_chinese

       #$($mode = "application.mode=dev") && echo "开发模式2" || echo "生产模式2"
      # return "未知模式"

}

function update()
{
  if [ -d "$BASE_DIRECTORY" ]
  then
       success
       cd $BASE_DIRECTORY && svn update || echo "更新失败" && exit 1
       RET_VAL=0
  else
       failure
       echo "$BASE_DIRECTORY不是目录"
       RET_VAL=1
  fi
  return $RETVAL

}
case $1 in
'start')
        #start的时候也调用restart 因为当play异常关闭后server.pid还在 此时若调用start命令并不能重启
        $0 restart
;;
'stop')
        stop
;;
'restart')
        stop
        sleep 1
        start
;;
'status')
        status
;;
'update')
        update
;;
*)
        echo "Usage: service sso  <start|stop|status|restart|update>"
;;
esac
exit 0

